# filesystem/etc/systemd/user/quadlet-sync.service
[Unit]
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/var/mnt/data/%N.env
ExecStart=/usr/local/bin/quadlet-sync.sh

[Install]
WantedBy=network-online.target