# filesystem/etc/systemd/user/quadlet-sync.timer
[Timer]
OnStartupSec=15s
OnUnitActiveSec=2m
Persistent=true
AccuracySet=1s

[Install]
WantedBy=timers.target