# filesystem/etc/systemd/user/quadlet-reload.path
[Path]
PathChanged=%h/git/infra

[Install]
WantedBy=default.target