# filesystem/etc/systemd/user/quadlet-sync.timer
[Timer]
OnBootSec=2m
OnUnitActiveSec=2m

[Install]
WantedBy=timers.target