# systemd/backup.service
[Unit]
Description=Backup

[Service]
ExecStart=%h/backup-podman.sh